# cron-jobs/backup-s3.cron
#!/bin/bash

/usr/bin/s3cmd sync /var/www/html/current/shared/db_backups s3://alchemy-external-backup/
/usr/bin/s3cmd sync /var/www/html/current/content/uploads s3://alchemy-external-backup/
